# video_collector.service
[Unit]
Description=Video Collector Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/video_collector
ExecStart=%h/video_collector/run.sh
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
